# com.brainwav.cortexdx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identification -->
    <key>Label</key>
    <string>com.brainwav.cortexdx</string>

    <!-- Program Configuration -->
    <key>ProgramArguments</key>
    <array>
        <string>__PNPM_BIN__</string>
        <string>--dir</string>
        <string>__PACKAGE_DIR__</string>
        <string>run</string>
        <string>server</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>__PORT__</string>
        <key>HOST</key>
        <string>__HOST__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <!-- Execution Settings -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/var/log/cortexdx.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/cortexdx.error.log</string>

    <!-- Process Management -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>ExitTimeOut</key>
    <integer>30</integer>

    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>ResidentSetSize</key>
        <integer>1342177280</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
        <key>ResidentSetSize</key>
        <integer>1610612736</integer>
    </dict>

    <!-- User Context -->
    <key>UserName</key>
    <string>__USER_NAME__</string>

    <key>GroupName</key>
    <string>__GROUP_NAME__</string>
</dict>
</plist>
