# com.brainwav.cortexdx.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brainwav.cortexdx.healthcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/scripts/monitoring/health-check.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>CORTEXDX_HEALTH_LOG_DIR</key>
        <string>__LOG_DIR__</string>
        <key>CORTEXDX_HEALTH_URLS</key>
        <string>http://127.0.0.1:5001/health https://cortexdx-mcp.brainwav.io/health</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/healthcheck.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/healthcheck.launchd.error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
